# systemd/user/hyprland-onboot.service
[Unit]
Description=Start Hyprland on boot

[Service]
Type=simple
ExecStart=/usr/bin/Hyprland

[Install]
WantedBy=default.target
